# docker/etc/supervisord.d/processes.ini
[program:smokeping]
command=/opt/smokeping/bin/smokeping --nodaemon --logfile=/dev/stdout
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
